# packages/worker/install/com.paperclip.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.paperclip.worker</string>

  <!-- No arguments needed — worker starts in beacon-only mode,
       auto-detects installed CLIs, and generates a pairing key.
       Pair from the Paperclip UI using the key in ~/.paperclip/worker-key -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/usr/local/lib/node_modules/@paperclipai/worker/dist/main.js</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/paperclip-worker.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/paperclip-worker.err</string>

  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
